# supervisord.d/mongodb.ini
[program:mongodb]
command=mongod --dbpath /mnt/project/db/ --logpath /mnt/project/log/mongodb.log --logappend --bind_ip '127.0.0.1' --journal ; the program (relative uses PATH, can take args)
stdout_logfile=/mnt/project/log/mongodb.stdout.log ; stdout log path, NONE for none; default AUTO
stderr_logfile=/mnt/project/log/mongodb.stderr.log ; stderr log path, NONE for none; default AUTO
