# supervisord.d/nodejs.ini
[program:nodejs]
command=node code/app          ; the program (relative uses PATH, can take args)
directory=/mnt/project/src     ; directory to cwd to before exec (def no cwd)
stdout_logfile=/mnt/project/log/nodejs.stdout.log ; stdout log path, NONE for none; default AUTO
stderr_logfile=/mnt/project/log/nodejs.stderr.log ; stderr log path, NONE for none; default AUTO
